# deploy/rfm-server.service
[Unit]
Description=Remote File Manager Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=rfm
Group=rfm
WorkingDirectory=/opt/remote-file-manager
ExecStart=/opt/remote-file-manager/server
Restart=always
RestartSec=5

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/remote-file-manager/data

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment="CONFIG_PATH=/opt/remote-file-manager/config.yaml"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rfm-server

[Install]
WantedBy=multi-user.target
